# RomanUApkg/DEBIAN/postinst
#!/bin/bash

set -e

systemctl daemon-reload

# Enable systemd service

systemctl enable skr_ip.service
systemctl start  skr_ip.service
